# build/src/cron.d/qsiphon
*/5 * * * * root /usr/local/bin/siphoncron -N > /dev/null
